# infra/systemd/user/tb1-robot-agent.service
[Unit]
Description=TB1 fleet robot status agent
After=tb1-network-ready.service tb1-bringup.service
Wants=tb1-network-ready.service tb1-bringup.service

[Service]
Type=simple
Environment=ROS_DISTRO=humble
Environment=ROS_DOMAIN_ID=42
Environment=ROS_LOCALHOST_ONLY=0
Environment=RMW_IMPLEMENTATION=rmw_cyclonedds_cpp
ExecStart=/usr/bin/bash -lc 'source /opt/ros/humble/setup.bash && source %h/turtlebot-fleet-ops/install/setup.bash && exec ros2 launch robot_agent robot_agent.launch.py'
Restart=on-failure
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=15

[Install]
WantedBy=default.target
