# infra/systemd/user/tb1-navigation.service
[Unit]
Description=TB1 local AMCL, Nav2, fleet lease agent, and motion arbiter
After=tb1-bringup.service tb1-safety-watchdog.service tb1-robot-agent.service
Wants=tb1-bringup.service tb1-safety-watchdog.service tb1-robot-agent.service
Conflicts=tb1-mapping.service

[Service]
Type=simple
Environment=ROS_DISTRO=humble
Environment=ROS_DOMAIN_ID=42
ExecStartPre=/usr/bin/test -r %h/.local/share/turtlebot-fleet-ops/maps/tb1/map.yaml
ExecStart=/usr/bin/bash %h/turtlebot-fleet-ops/infra/navigation/start-tb1-navigation.sh
ExecStartPost=/usr/bin/systemctl --user try-restart tb1-zenoh-bridge.service
Restart=always
RestartSec=3
KillSignal=SIGINT
TimeoutStopSec=20

[Install]
WantedBy=default.target
